# dwupdate.service
[Unit]
Description=DWM statusbar

[Service]
Type=simple
ExecStart=/usr/bin/dwupdate

[Install]
WantedBy=graphical.target
